# systemd/xnetvn_monitord.service
[Unit]
Description=xNetVN Server Monitor Daemon
Documentation=https://github.com/xntmkt/servers
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/xnetvn_monitord
ExecStart=/usr/bin/python3 -m xnetvn_monitord.daemon /opt/xnetvn_monitord/config/main.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xnetvn_monitord

# Security settings
PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
